# scripts/voice-bridge.service
[Unit]
Description=Voice Bridge - Offline Voice Assistant Service
After=network.target

[Service]
Type=simple
User=%I
WorkingDirectory=/usr/local/lib/voice-bridge
ExecStart=/usr/bin/python3 /usr/local/lib/voice-bridge/start_adapters.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
